# systemd/.config/systemd/user/compton.service
[Unit]
Description=Start Compton Window Compositor
After=monitor.service

[Service]
Type=simple
ExecStart=compton

[Install]
WantedBy=init.service
